# user/systemd-user-laptop/.config/systemd/user/dunst.service.d/override.conf
[Unit]
Require=windows-manager.service
After=windows-manager.service

[Install]
WantedBy=
WantedBy=windows-manager.service
