# user/systemd-user-laptop/.config/systemd/user/i3.service
[Unit]
Description=i3 window manager

[Service]
ExecStart=/usr/bin/i3 -V

[Install]
Alias=windows-manager.service
