# user/systemd-user-laptop/.config/systemd/user/redshift.service.d/override.conf
[Install]
WantedBy=
WantedBy=windows-manager.service
